# bin/meshtasticd.service
[unit]
description=Meshtastic Native Daemon

[Service]
Type=simple
ExecStart=/usr/sbin/meshtasticd

[Install]
WantedBy=multi-user.target
